# packages/bng-stack-runtime/etc/systemd/system/vpp.service.d/stderr-suppress.conf
[Service]
# Suppress VPP stderr to prevent log floods from stale API client warnings
StandardError=null
